# streamserver.service
[Unit]
Description=StreambotV2 Stream Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/StreambotV2
ExecStart=/root/StreambotV2/venv/bin/python stream_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
